# root/etc/services.d/srv/run
#!/usr/bin/with-contenv sh
su user -c /server/server