# root/etc/services.d/srv/run
#!/usr/bin/with-contenv sh
# copy to allow for dev to overwrite the server
cp /server/server /tmp/server
su user -c /tmp/server