# root/etc/services.d/postgrest/run
#!/usr/bin/with-contenv sh

su node -c '/usr/local/bin/postgrest /etc/postgrest.conf'

