# root/etc/services.d/postgrest/run
#!/usr/bin/with-contenv sh

/usr/local/bin/postgrest /etc/postgrest.conf
