# root/etc/services.d/nginx/run
#!/usr/bin/with-contenv sh

nginx -g 'daemon off; include /etc/nginx/conf.d/*.conf;'
